# state/apache2/files/vnstat.cron
*/5 * * * * root /bin/sh /opt/fffd-keyupload/script/vnstat.sh
