# state/bind/files/fffd-dns-nodes.cron
#!/bin/sh

# THIS FILE IS CONTROLLED BY SALTSTACK!
#
# Generate the dynamic DNS zone nodes.fffd.eu

if [ ! -d /var/lib/bind/zones ]; then
  mkdir -p /var/lib/bind/zones
fi

/opt/fffd-dns/dns-nodenames.py >/var/lib/bind/zones/db.nodes.fffd.eu
/usr/sbin/rndc reload &>/dev/null

exit 0
